# systemd/meross-daemon.service
[Unit]
Description=Meross MSS310 power monitoring daemon
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=pi
ExecStart=/usr/bin/python3 /home/pi/Desktop/meross_daemon.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
